# deploy/ironbuddy-supervisor.conf
; IronBuddy V2.5 — Supervisor 进程守护配置
; 用法: supervisord -n -c /home/toybrick/deploy/ironbuddy-supervisor.conf
; 安装: sudo apt-get install supervisor

[unix_http_server]
file=/tmp/ironbuddy-supervisor.sock

[supervisord]
logfile=/tmp/supervisord-ironbuddy.log
pidfile=/tmp/supervisord-ironbuddy.pid
childlogdir=/tmp

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/ironbuddy-supervisor.sock


; ===== NPU 姿态推理引擎 =====
[program:npu]
command=/home/toybrick/yolo_test/build/main
directory=/home/toybrick/yolo_test/build
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/tmp/ironbuddy-npu.log
stdout_logfile_maxbytes=2MB
stderr_redirect_stderr=true
priority=100

; ===== Flask 推流中台 =====
[program:streamer]
command=python3 /home/toybrick/streamer_app.py
directory=/home/toybrick
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/tmp/ironbuddy-streamer.log
stdout_logfile_maxbytes=2MB
stderr_redirect_stderr=true
priority=200

; ===== FSM 主循环 + V2 管线 =====
[program:fsm]
command=python3 /home/toybrick/hardware_engine/main_claw_loop.py
directory=/home/toybrick/hardware_engine
environment=OPENCLAW_URL="ws://127.0.0.1:18789"
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/tmp/ironbuddy-fsm.log
stdout_logfile_maxbytes=2MB
stderr_redirect_stderr=true
priority=300

; ===== 语音守护进程 =====
[program:voice]
command=python3 /home/toybrick/hardware_engine/voice_daemon.py
directory=/home/toybrick/hardware_engine
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/tmp/ironbuddy-voice.log
stdout_logfile_maxbytes=2MB
stderr_redirect_stderr=true
priority=400

; ===== TTS 守护 =====
[program:tts]
command=bash /home/toybrick/hardware_engine/peripherals/tts_daemon.sh
directory=/home/toybrick/hardware_engine/peripherals
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/tmp/ironbuddy-tts.log
stdout_logfile_maxbytes=2MB
stderr_redirect_stderr=true
priority=500

; ===== 外设旁路监听 =====
[program:peripheral]
command=bash /home/toybrick/hardware_engine/peripherals/peripheral_daemon.sh
directory=/home/toybrick/hardware_engine/peripherals
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/tmp/ironbuddy-peripheral.log
stdout_logfile_maxbytes=2MB
stderr_redirect_stderr=true
priority=600


; ===== 进程组 =====
[group:ironbuddy]
programs=npu,streamer,fsm,voice,tts,peripheral
priority=999
